# rules/packlim/etc/profile.d/packlim
export REPO="http://rlsd3.dimakrasner.com/repo/$(uname -m)"
