# rules/loksh/etc/profile.d/loksh
export PS1='[\u@\w]\$ '
export HISTFILE="$HOME/.ksh_history"
alias ll="ls -la"
