# rules/packlad/etc/profile.d/packlad
export REPO="http://rlsd2.dimakrasner.com/repo"
